# source/strafer_ros/99-strafer.rules
# Udev rules for Strafer robot persistent device names.
#
# Two RoboClaws (Basicmicro) share identical vendor/product/serial attributes.
# The driver auto-detects which is front vs rear by probing each port's packet
# serial address (0x80 = front, 0x81 = rear) at startup.
#
# This rule creates generic symlinks /dev/roboclaw0, /dev/roboclaw1 for any
# Basicmicro RoboClaw regardless of which USB port it's plugged into.
# The driver's auto-detect scans /dev/roboclaw* to find candidates.
#
# Install:  make udev   (or manually):
#   sudo cp 99-strafer.rules /etc/udev/rules.d/
#   sudo udevadm control --reload-rules
#   sudo udevadm trigger
#
# Verify:
#   ls -la /dev/roboclaw*

# Any Basicmicro RoboClaw -> /dev/roboclaw<N>
SUBSYSTEM=="tty", ATTRS{idVendor}=="03eb", ATTRS{idProduct}=="2404", SYMLINK+="roboclaw%n", MODE="0666"

# Intel RealSense D555 HID interface (IMU: BMI055 accel + gyro)
# Custom-built hid-sensor-hub + IIO modules expose the IMU as IIO devices.
# Grant rw access to /dev/iio:device* and fix sysfs scan_elements permissions
# so librealsense can read IMU data without root.
SUBSYSTEM=="iio", KERNEL=="iio:device*", MODE="0666", RUN+="/bin/sh -c 'find /sys%p -type f -exec chmod a+rw {} + 2>/dev/null'"

# Also keep hidraw accessible for the D555 HID interface
KERNEL=="hidraw*", ATTRS{idVendor}=="8086", ATTRS{idProduct}=="0b56", MODE="0666"
